# cmd/installer/include-linux/eco.service
[Unit]
Description=Decred Eco system service

[Service]
ExecStart=<ExecStart>

[Install]
WantedBy=multi-user.target