# installers/debian/eco.service
[Unit]
Description=Decred Eco system service

[Service]
ExecStart=/usr/local/bin/ecoservice

[Install]
WantedBy=multi-user.target